# overlay/etc/periodic/daily/update-certs
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Run SSL update executable and log output to /ssl directory
#======================================================================================================================

ssl-update > ${SSL}/update.log
