# overlay/etc/periodic/daily/update-certs
#!/bin/bash

set -euo pipefail
ssl-update
