# overlay/etc/periodic/weekly/update-certs
#!/usr/bin/with-contenv bash
ssl-update > ${SSL}/update.log
