# overlay/etc/periodic/weekly/update-certs
#!/bin/bash
ssl-update > /ssl/update.log
