# overlay/etc/s6-overlay/s6-rc.d/ssl-auto-request/run
#!/usr/bin/nu

use bf
use bf/nginx/proxy auto
bf env load -x ssl-auto-request

# Request SSL certificates and then disable the service
def main [...args] {
    if (auto is_enabled) {
        # get the pid of the nginx process - if the pid is empty, nginx is not running
        let pid = { ^pidof nginx } | bf handle
        if $pid == "" {
            # wait 2s before exiting the service - S6 will keep restarting it until Nginx comes online
            # on first run, it will disable this upgrade service itself
            let sleep_for = 2sec
            bf write debug $"Waiting ($sleep_for) for Nginx to come online."
            sleep $sleep_for
        } else {
            # request SSL certificates for all configured domains

            # disable the auto request service
            auto disable_svc
        }
    } else {
        # disable the auto request service
        auto disable_svc
    }
}
